# macOS-Samples/Scripts/Legacy App Catalog Webclip For All Docks/payload/Library/LaunchAgents/com.vmware.catalogtodock.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.vmware.catalogtodock</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/python</string>
      <string>/Users/Shared/VMware/catalogtodock.py</string>
    </array>
    <key>OnDemand</key>
    <true/>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/Shared/VMware/catalogtodock.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/Shared/VMware/catalogtodock.log</string>
  </dict>
</plist>
